# ansible-practice-realtime/server_cfg/web-server.service
[Unit]
Description=Web-Server
[Service]
Type=idle
WorkingDirectory=/home/armitage
ExecStart=python3 app.py
[Install]
WantedBy=multi-user.target

